# systemd/openclaw-node@.service
[Unit]
Description=OpenClaw Node Host (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/openclaws/pkg/bin/run_node_instance.sh %i
Restart=on-failure
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target

